# private_dot_config/systemd/user/claude-code-update.service
[Unit]
Description=Claude Code auto-update check

[Service]
Type=oneshot
Environment=HOME=%h
EnvironmentFile=-%h/.config/ai-updaters/env
ExecStart=%h/.local/share/ai-updaters/bin/upgrade-claude-code
TimeoutStartSec=300
Restart=no
