# private_dot_config/systemd/user/claude-code-update.timer
[Unit]
Description=Every 5 minutes Claude Code update check

[Timer]
OnStartupSec=5min
OnUnitActiveSec=5min
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
